# ref/desktop-arch/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin golly --noclear %I $TERM
TTYVTDisallocate=no
