# ref/xps-arch/etc/systemd/system/graphical.target.wants/bumblebeed.service
[Unit]
Description=Bumblebee C Daemon

[Service]
Type=simple
CPUSchedulingPolicy=idle
ExecStart=/usr/bin/bumblebeed
Restart=always
RestartSec=60
StandardOutput=kmsg

[Install]
WantedBy=graphical.target

