# ref/xps-arch/home/.config/systemd/user/default.target.wants/bloop.service
[Unit]
Description=Bloop Scala build server

[Service]
ExecStart=/home/golly/.bloop/bloop server
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bloop

[Install]
WantedBy=default.target
