# ref/filelist
~/.config/autostart/
~/.config/compiz/compizconfig/Default.ini
~/.config/monitors.xml
~/.config/smplayer/smplayer.ini
~/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-keyboard-shortcuts.xml
~/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-session.xml
~/.config/xfce4/xfconf/xfce-perchannel-xml/xsettings.xml
~/.conkyrc
~/.colordiffrc
~/.cups/lpoptions
~/.m2/settings.xml
~/.mpv
~/.xinitrc
/etc/adobe/mms.cfg
/etc/conf.d/
/etc/default/tlp
/etc/fstab
/etc/hosts
/etc/httpd/conf/httpd.conf
/etc/httpd/conf/mime.types
/etc/i8kutils/i8kmon.conf
/etc/inittab
/etc/inputrc
/etc/locale.conf
/etc/makepkg.conf
/etc/mkinitcpio.conf
/etc/modprobe.d/
/etc/modules-load.d/
/etc/ntp.conf
/etc/php/php.ini
/etc/postfix/main.cf
/etc/profile.d/disable_beep.sh
/etc/rc.conf
/etc/rc.local
/etc/resolv.conf
/etc/sysctl.conf
/etc/sysctl.d/
/etc/systemd/system/
/etc/udev/rules.d/
/etc/vconsole.conf
/etc/X11/xorg.conf
/etc/X11/xorg.conf.d/10-monitor.conf
/etc/X11/xorg.conf.d/20-intel.conf
/etc/yaourtrc
