# ref/xps-arch/etc/systemd/system/multi-user.target.wants/powertop.service
[Unit]
Description=Powertop tunings

[Service]
ExecStart=/usr/bin/powertop --auto-tune
Type=oneshot
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
