# ref/desktop-arch/home/.config/autostart/keybase_autostart.desktop
[Desktop Entry]
Name=Keybase
Comment=Keybase Filesystem Service and GUI
Type=Application
Exec=env KEYBASE_START_UI=hideWindow run_keybase
Hidden=true

